# deploy/systemd/rootstone-dashboard.service
[Unit]
Description=Rootstone Dashboard Next.js App
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/rootstone-dashboard
Environment=HOME=/home/ubuntu
Environment=NODE_ENV=production
Environment=HOSTNAME=0.0.0.0
Environment=PORT=3000
EnvironmentFile=/home/ubuntu/rootstone-dashboard/.env
ExecStart=/usr/bin/env bash -lc '. /home/ubuntu/.nvm/nvm.sh && exec node .next/standalone/server.js'
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
