# deploy/systemd/rootstone-dashboard-daily-nav.service
[Unit]
Description=Rootstone Dashboard daily NAV cron
After=rootstone-dashboard.service

[Service]
Type=oneshot
User=ubuntu
WorkingDirectory=/home/ubuntu/rootstone-dashboard
Environment=HOME=/home/ubuntu
EnvironmentFile=/home/ubuntu/rootstone-dashboard/.env
ExecStart=/usr/bin/env bash /home/ubuntu/rootstone-dashboard/deploy/scripts/run-cron-route.sh daily-nav
